# tests/profile.d/c7epas.profile
export CLNAME="c7epas"
export DBTYPE="EPAS"
export DBVERSION="15"
export DOCKERI="centos:7"
