# tests/profile.d/c7pg.profile
export CLNAME="c7pg"
export DBTYPE="PG"
export DBVERSION="16"
export DOCKERI="centos:7"
