# tests/profile.d/c8pg.profile
export CLNAME="c8pg"
export DBTYPE="PG"
export DBVERSION="14"
export DOCKERI="centos:8"
