# tests/profile.d/d10epas.profile
export CLNAME="d10epas"
export DBTYPE="EPAS"
export DBVERSION="13"
export DOCKERI="debian:10"
