# tests/profile.d/d10pg.profile
export CLNAME="d10pg"
export DBTYPE="PG"
export DBVERSION="13"
export DOCKERI="debian:10"
