# tests/profile.d/d11epas.profile
export CLNAME="d11epas"
export DBTYPE="EPAS"
export DBVERSION="15"
export DOCKERI="debian:11"
