# tests/profile.d/d9epas.profile
export CLNAME="d9epas"
export DBTYPE="EPAS"
export DBVERSION="14"
export DOCKERI="debian:9"
