# tests/profile.d/d9pg.profile
export CLNAME="d9pg"
export DBTYPE="PG"
export DBVERSION="15"
export DOCKERI="debian:9"
