# tests/profile.d/ro8epas.profile
export CLNAME="ro8epas"
export DBTYPE="EPAS"
export DBVERSION="15"
export DOCKERI="rockylinux:8"
