# tests/profile.d/ro8pg.profile
export CLNAME="ro8pg"
export DBTYPE="PG"
export DBVERSION="15"
export DOCKERI="rockylinux:8"
