# tests/profile.d/ro9pg.profile
export CLNAME="ro9pg"
export DBTYPE="PG"
export DBVERSION="16"
export DOCKERI="rockylinux:9"
