# tests/profile.d/u18epas.profile
export CLNAME="u18epas"
export DBTYPE="EPAS"
export DBVERSION="14"
export DOCKERI="ubuntu:18.04"
