# tests/profile.d/u20epas.profile
export CLNAME="u20epas"
export DBTYPE="EPAS"
export DBVERSION="15"
export DOCKERI="ubuntu:20.04"
