# tests/profile.d/u20pg.profile
export CLNAME="u20pg"
export DBTYPE="PG"
export DBVERSION="16"
export DOCKERI="ubuntu:20.04"
