# tests/profile.d/u22pg.profile
export CLNAME="u22pg"
export DBTYPE="PG"
export DBVERSION="16"
export DOCKERI="ubuntu:22.04"
